# doc/utils/crontab
BIO_PATH=~/Projects/biopython/doc/utils

30 2 * * * /your/command
